# deploy/bglib.service
[Unit]
Description=BgLib board game library (Next.js)
After=network.target

[Service]
Type=simple
User=bglib
Group=bglib
WorkingDirectory=/opt/bglib
EnvironmentFile=/opt/bglib/.env.local
Environment=NODE_ENV=production
Environment=BIND_HOST=0.0.0.0
Environment=PORT=3000
ExecStart=/bin/bash /opt/bglib/scripts/start-production.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
