# deploy/bglib-cron.timer
[Unit]
Description=Daily BgLib loan reminders at 09:00 UTC

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
